# to_install/etc/systemd/system/setup-and-update.target
[Unit]
Description=Activated set up and update service
Wants=setup-and-update.service
